# openzephir-bot/systemd/openzephir-bot.service
[Unit]
Description=OpenZephir Discord Bot
After=network.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/openzephir-bot
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
